# deploy/nemh-api.service
# NEMH 后端 API — systemd 单元
#
# 安装（在服务器上，按需改路径与 User）：
#   sudo bash deploy/install-nemh-api-systemd.sh
# 或手动：
#   sudo mkdir -p /var/log/nemh-api && sudo chown ubuntu:ubuntu /var/log/nemh-api
#   sudo cp deploy/nemh-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nemh-api
#
# 日志：
#   应用输出：tail -f /var/log/nemh-api/out.log  与  tail -f /var/log/nemh-api/err.log
#   服务启停：journalctl -u nemh-api -f
#
# 若改用 journald 承接应用 stdout（不用文件），将 StandardOutput/StandardError
# 改为 journal，并注释掉 append 两行，同时可设 SyslogIdentifier=nemh-api

[Unit]
Description=NEMH API (Express, PORT=8003)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# 修改为服务器上本仓库的 server 目录绝对路径
WorkingDirectory=/home/ubuntu/var/www/nemh-app/NEMH/server

Environment=NODE_ENV=production
Environment=PORT=8003

# 使用系统 Node（与 apt / NodeSource 一致）；若仅用 nvm，请改为该 node 绝对路径
ExecStart=/usr/bin/node src/index.js

Restart=on-failure
RestartSec=5
TimeoutStopSec=20

# 应用日志（需目录存在且 User 可写）
StandardOutput=append:/var/log/nemh-api/out.log
StandardError=append:/var/log/nemh-api/err.log

[Install]
WantedBy=multi-user.target
