# deploy/nemh-api.service
# NEMH 后端 API — systemd 单元
#
# 安装：在仓库根目录执行
#   sudo bash deploy/install-nemh-api-systemd.sh [REPO_ROOT]
#
# 日志（相对仓库根目录，占位 __REPO_ROOT__ 由安装脚本替换为绝对路径）：
#   logs/api.log     — Node 标准输出 + 标准错误（应用日志）
#   logs/service.log — 服务启动/停止时间戳（ExecStartPre / ExecStopPost）
# 详细 systemd 消息仍可用：journalctl -u nemh-api -f
#
# 若改用 journal 承接应用输出：将 StandardOutput/StandardError 改为 journal，
# 并注释掉下面两行 append。

[Unit]
Description=NEMH API (Express, PORT=8003)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# 服务器上本仓库的 server 目录（占位由 install 脚本替换）
WorkingDirectory=__REPO_ROOT__/server

Environment=NODE_ENV=production
Environment=PORT=8003

# Project2（PD_max）用户库与 SSO — 生产启用时取消注释并核对密钥
Environment=PD2_AUTH_ENABLED=1
Environment=PD2_MYSQL_HOST=127.0.0.1
Environment=PD2_MYSQL_PORT=3306
Environment=PD2_MYSQL_USER=root
Environment=PD2_MYSQL_PASSWORD=请改为服务器上 PD_max/.env 中的 MYSQL_PASSWORD
Environment=PD2_MYSQL_DATABASE=PD_max
Environment=PD2_JWT_SECRET=change_this_to_a_strong_random_secret
Environment=PD2_API_BASE_URL=http://127.0.0.1:8001
Environment=DISABLE_PUBLIC_REGISTER=1

ExecStartPre=/bin/sh -c 'mkdir -p __REPO_ROOT__/logs && echo "$$(date -Is) nemh-api starting" >> __REPO_ROOT__/logs/service.log'

# 使用系统 Node；若仅用 nvm，改为该 node 绝对路径
ExecStart=/usr/bin/node src/index.js

ExecStopPost=/bin/sh -c 'echo "$$(date -Is) nemh-api stopped" >> __REPO_ROOT__/logs/service.log'

Restart=on-failure
RestartSec=5
TimeoutStopSec=20

# 应用日志（stdout + stderr 合并到 api.log）
StandardOutput=append:__REPO_ROOT__/logs/api.log
StandardError=append:__REPO_ROOT__/logs/api.log

[Install]
WantedBy=multi-user.target
